# supervisor.conf
[supervisord]
loglevel=debug
pidfile=supervisord.pid
nodaemon=true
childlogdir=.

[program:client]
process_name=client-%(process_num)d
command=python zmq-client.py
numprocs=4
[program:server]
process_name=server-%(process_num)d
command=python zmq-server.py
numprocs=2
[program:queue]
command=python zmq-queue.py
[program:web]
command=python app.py
